# Dockerfiles/system/s6/services.d/rsyslog/run
#!/bin/bash
source /etc/sysconfig/rsyslog
/usr/sbin/rsyslogd -n