# Dockerfiles/backend/mysql/s6/services.d/mysql/run
#!/bin/bash
## wait mysql start and init mysql

/usr/sbin/mysqld --user=root
