# Dockerfiles/backend/mysql/s6/services.d/mysql/run
#!/bin/bash
/usr/sbin/mysqld --user=root
